# deploy/systemd/ugnosd.service
# systemd unit for ugnosd. Install: copy to /etc/systemd/system/, then:
#   systemctl daemon-reload
#   systemctl enable --now ugnosd
# Adjust paths and User/Group for your deployment.

[Unit]
Description=UGNOS time-series database daemon
Documentation=https://github.com/bogwi/ugnos
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ugnos
Group=ugnos
ExecStart=/usr/local/bin/ugnosd --data-dir /var/lib/ugnos --http-bind 0.0.0.0:8080
# Or with config file:
# ExecStart=/usr/local/bin/ugnosd --config /etc/ugnosd.toml
ExecReload=/bin/kill -TERM $MAINPID
TimeoutStopSec=35
Restart=on-failure
RestartSec=5

# Security hardening (optional)
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/lib/ugnos
# If using config in /etc:
# ReadOnlyPaths=/etc/ugnosd.toml

[Install]
WantedBy=multi-user.target
